# UNIX/etc/xdg/autostart/neubot.desktop
[Desktop Entry]
Type=Application
Exec=@BINDIR@/neubot notifier
Hidden=False
NoDisplay=False
X-GNOME-Autostart-enabled=true
Name=Neubot notifier
