# M-Lab/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

#
# Modified /etc/rc.local for Neubot
# Installed by M-Lab/install.sh
#

touch /var/lock/subsys/local
/home/mlab_neubot/neubot/M-Lab/start.sh
